# crontab
# 每 30 分钟执行一次爬虫脚本
*/30 * * * * /usr/local/bin/php /app/task.php >> /app/cron.log 2>&1
